# contrib/castblock.service
[Unit]
Description=CastBlock Service
After=network.target

[Service]
ExecStart=/usr/bin/castblock
